# magic-serve.conf
[program:magic-serve]
command=/usr/local/bin/gunicorn wsgi --bind localhost:2000 --workers 3 --worker-class gevent
directory=/var/www/magic-serve
autostart=true
autorestart=true
